# scraperapp/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:dockerlogging]
command=/opt/log.sh
startsecs = 0
autorestart = false
startretries = 1