# britishcolumbia-spider/supervisord.conf
[supervisord]
nodaemon=true

[program:dockerlogging]
command=/opt/log.sh
startsecs = 0
autorestart = false
startretries = 1